# DDIT/DDIT_supervisor/nginx.conf
[program:nginx]
command = nginx -g 'daemon off;'
startsecs=0
autostart=true
autorestart=true
#stopasgroup=true
#killasgroup=true`

